# systemd_files/Backup_Int_2_Ext.service
# A systemd service to trigger ripping of blu rays, this script is called by a udev rule
# this script should be stored in /etc/sytemd/system with permissions 444

[Unit]
Description=Auto Backup Internal HDD - External USB

[Service]
Type=oneshot
ExecStart=/usr/local/bin/AutoBackupInt2Ext.sh

[Install]
WantedBy=multi-user.target
